# systemd/wg-all-restart.path
[Unit]
Description=Watch /etc/wireguard/wg-all.conf for changes

[Path]
PathModified=/etc/wireguard/wg-all.conf

[Install]
WantedBy=multi-user.target
